# phusion-cloud9/etc/service/cloud9/run
#!/bin/sh
export HOME="/root"
exec /usr/bin/nodejs /cloud9/server.js -w /var/www/html --listen 0.0.0.0 -p 80 -a : 
